# nvidia_exporter.conf
# Run nvidia_exporter


start on startup

script
    /usr/bin/nvidia_exporter
end script
